# campaign-db/systemd/campaign-db-restore-test.timer
[Unit]
Description=Run the campaign-settings restore test weekly (Sun 04:30 local)

[Timer]
# After the 03:00 nightly backup, so the newest dump is that night's.
OnCalendar=Sun *-*-* 04:30:00
RandomizedDelaySec=10min
Persistent=true

[Install]
WantedBy=timers.target
